# com.user.remountnetworkdrive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.remountnetworkdrive</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/path/to/media_nas_connect.sh</string>
    </array>

    <!-- Run the script at load -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Only restart the script if it exits with a non-zero status (failure).
         SuccessfulExit: false means "do NOT keep alive after a clean exit". -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Minimum seconds between restarts to prevent rapid respawn loops -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Start the script when volumes are mounted -->
    <key>StartOnMount</key>
    <true/>
</dict>
</plist>
